# apps/edge/systemd/fartola.service
# Example systemd user unit for the fartOLa edge bridge.
# Locked by:
#   .planning/phases/01-single-laptop-training-mvp/01-18-PLAN.md task 2
#   .planning/REQUIREMENTS.md REQ-OPS-002 (bridge survives restart with
#   zero data loss; SQLite WAL + cold-start migrator from plan 02
#   provide the data-loss guarantee, this unit drives the restart).
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp $(npm prefix -g)/lib/node_modules/@fartola/edge/systemd/fartola.service \
#      ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now fartola
#   systemctl --user status fartola
#
# `%h` expands to the invoking user's home directory; using the user-scope
# manager (systemctl --user) means we don't need root and the unit is
# isolated from system services. Set FARTOLA_DEV=1 in an override file
# (systemctl --user edit fartola) ONLY on bench / single-operator laptops —
# the dev gate exposes /api/__dev/* + /api/__admin/* (T-DEV-ENDPOINT,
# T-ADMIN-DEV-GATE) which Phase 2 will replace with admin-token auth.

[Unit]
Description=fartOLa training MVP edge bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/fartola \
  --port 3000 \
  --bind-host 127.0.0.1 \
  --db-path %h/.local/share/fartola/fartola.db \
  --backup-dir %h/.local/share/fartola/backups \
  --retention-days 30
Restart=on-failure
RestartSec=5s

# Hardening — operator unit, no root needed, but tighten what we can.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/fartola
PrivateTmp=true

# Logging routed to the user journal; tail with `journalctl --user -u fartola -f`.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
